# etc/profile.d/proxy.sh
export http_proxy=http://172.16.30.20:8080
export https_proxy=http://172.16.30.20:8080
